# src/systemd/topsocialtwitter.service
[Unit]
After=arangodb3.service
Description=Topsocial Twitter Service

[Service]
Type=idle
User=tscl
ExecStart=/usr/bin/python3.6 -m core.twitter_top_social
WorkingDirectory=/home/tscl/gopath/src/topsocial/src/

[Install]
WantedBy=multi-user.target


